# ops/cerbot_renew.cron
#!/bin/sh

set -e

# Monthly job to renew Lets Encrypt certificate
certbot renew --post-hook "service restart nginx"
