# ops/thekevincrane_upstart.conf
start on [2345]
stop on [06]

env APP_ENV=prod

script
    exec uwsgi --emperor /etc/uwsgi/apps-enabled
end script